# chillhub-host/chillhub-init.service
[Unit]
Description=ChillHub Init Service
Documentation=
After=

[Service]
ExecStart=/usr/bin/chillhub init

[Install]
WantedBy=multi-user.target
